# scripts/smart-queue.service
[Unit]
Description=Smart Queue Daemon - Character Library Crawler
After=network.target
Wants=network.target
# Only start if daemon script exists
ConditionPathExists=/home/pedro/projetos/char-lib/scripts/smart-queue-daemon.js

[Service]
Type=simple
User=smartqueue
Group=smartqueue
WorkingDirectory=/home/pedro/projetos/char-lib
ExecStart=/usr/bin/node ./scripts/smart-queue-daemon.js
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -USR1 $MAINPID
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=SMART_QUEUE_DATA_DIR=/home/pedro/projetos/char-lib/data
StandardOutput=journal
StandardError=journal
SyslogIdentifier=smart-queue

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
# Allow access to project files and config under /etc
ProtectHome=no
ReadWritePaths=/home/pedro/projetos/char-lib /etc/smart-queue /var/log /var/run/smart-queue

# Resource limits
MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target